# analytics/crontab
*/5 * * * * /usr/local/bin/python /app/analytics.py
